# base-docker/port_mapping.conf
[program:port-mapping]
directory=/tmp
command=/tmp/port_mapping.sh
redirect_stderr=true
priority=10
